# vespabase/src/vespa-configserver.service.in
[Unit]
Description=Vertical Search Platform Config Server

[Service]
Type=forking
PIDFile=@CMAKE_INSTALL_PREFIX@/var/run/configserver.pid
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/vespa-start-configserver.sh
ExecStop=@CMAKE_INSTALL_PREFIX@/bin/vespa-stop-configserver.sh

[Install]
WantedBy=multi-user.target