# vespabase/src/vespa-configserver.service.in
# Copyright Yahoo. Licensed under the terms of the Apache 2.0 license. See LICENSE in the project root.
[Unit]
Description=Vertical Search Platform Config Server
After=network.target

[Service]
Type=forking
PIDFile=@CMAKE_INSTALL_PREFIX@/var/run/configserver.pid
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/vespa-start-configserver
ExecStop=@CMAKE_INSTALL_PREFIX@/bin/vespa-stop-configserver

[Install]
WantedBy=multi-user.target
